# artemis.service
[Unit]
Description=Artemis Threat Hunting Platform (LAN-accessible on port 6969)
After=network.target

[Service]
Type=simple
User=administrator
WorkingDirectory=/home/user/Artemis
Environment="SPLUNK_HOST=10.25.11.86"
Environment="ARTEMIS_SYSTEMD=1"

# Load secrets from .env if it exists (the dash prefix means "don't fail if missing")
EnvironmentFile=-/home/user/Artemis/.env

ExecStart=/usr/local/bin/python3 artemis_server.py

# Kill the entire process group (server + all child pipelines/profilers)
KillMode=control-group
KillSignal=SIGTERM
TimeoutStopSec=15

Restart=always
RestartSec=10

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/home/user/Artemis
PrivateTmp=true

[Install]
WantedBy=multi-user.target
